# Firmware/systemd/mothbox-gpio.service
[Unit]
Description=Mothbox GPIO Daemon
After=network.target

[Service]
Type=notify
NotifyAccess=main
User=pi
Group=gpio
EnvironmentFile=/etc/mothbox/gpio-daemon.env
ExecStart=/usr/bin/python3 ${DAEMON_PATH}
Restart=on-failure
RestartSec=2
RuntimeDirectory=mothbox

[Install]
WantedBy=multi-user.target
