# Firmware/services/gps-exif-tagger.service
[Unit]
Description=Mothbox GPS EXIF Tagger
Documentation=https://github.com/Digital-Naturalism-Laboratories/Mothbox
After=network.target gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=pi
Group=pi

# Working directory (adjust based on installation type)
# Production: /opt/mothbox
# Legacy: /home/pi/Desktop/Mothbox/Firmware
# Development: /path/to/your/repo
WorkingDirectory=/opt/mothbox

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Main service command (watch mode with 10 second polling)
# Production installation:
ExecStart=/usr/bin/python3 /opt/mothbox/gps_exif_tagger.py --mode immediate --watch --interval 10 --verbose

# For legacy installation, use:
# ExecStart=/usr/bin/python3 /home/pi/Desktop/Mothbox/Firmware/gps_exif_tagger.py --mode immediate --watch --interval 10 --verbose

# Restart policy
Restart=on-failure
RestartSec=30s

# Logging (systemd journal)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gps-exif-tagger

# Resource limits (prevent runaway resource usage)
MemoryMax=256M
CPUQuota=25%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Filesystem protection
ProtectSystem=strict
ReadWritePaths=/var/lib/mothbox/photos
ReadOnlyPaths=/etc/mothbox

# Additional security layers
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Capability restrictions
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
